# deploy/llama.cpp.service.example
[Unit]
Description=llama.cpp translate server (MI50, Vulkan, gemma-4)
After=network-online.target
Wants=network-online.target

[Service]
ExecStartPre=/usr/bin/rocm-smi -d 0 --setpoweroverdrive 150
ExecStart=/opt/llama/llama.cpp/build-vulkan/bin/llama-server -hf unsloth/gemma-4-26B-A4B-it-qat-GGUF -ngl 99 -c 16384 --parallel 4 --n-predict 1024 --host 0.0.0.0 --port 8080
Restart=on-failure

[Install]
WantedBy=multi-user.target
